# skel/config/autostart/nitrogen.desktop
[Desktop Entry]
Exec=nitrogen --restore
Name=nitrogen
Type=Application
Version=1.0
