# skel/config/autostart/xfce-panel.desktop
[Desktop Entry]
Exec=xfce4-panel
Name=xfce-panel
Type=Application
Version=1.0
